# systemd/pga-betting-ai.timer
[Unit]
Description=PGA Betting AI weekly scan
Requires=pga-betting-ai.service

[Timer]
# Fire every Tuesday at 08:00 (day before most PGA Wednesday pro-ams)
OnCalendar=Tue *-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
